# etc/services.d/ovpn/run
#!/usr/bin/execlineb -P

cd /ovpn
openvpn --config client.conf --status status
